# plant-care-mcp.service
[Unit]
Description=Plant Care MCP Server
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/claude-code-the-gardener/app
Environment="PATH=/home/pi/.local/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/pi/claude-code-the-gardener/app/.env
ExecStart=/home/pi/.local/bin/uv run python run_http.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
