# static_site_generator/deploy/gardener-data-sync.timer
[Unit]
Description=Run Gardener Data Sync every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
AccuracySec=1min

[Install]
WantedBy=timers.target
