# static_site_generator/deploy/gardener-r2-sync.timer
[Unit]
Description=Run Gardener R2 Sync every 15 minutes
Documentation=https://github.com/mnbf9rca/claude-code-the-gardener

[Timer]
# Start 5 minutes after boot
OnBootSec=5min

# Run every 15 minutes
OnUnitActiveSec=15min

# If Pi was off, run missed syncs when it comes back online
Persistent=true

[Install]
WantedBy=timers.target
