# app/deploy/plant-care-mcp.service
[Unit]
Description=Plant Care MCP Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mcpserver
Group=mcpserver
WorkingDirectory=/home/mcpserver/plant-care-app
Environment="PATH=/home/mcpserver/.local/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/mcpserver/.env
ExecStart=/home/mcpserver/.local/bin/uv run python run_http.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Treat SIGTERM (exit code 143) as success - this is a graceful shutdown
SuccessExitStatus=143

# Graceful shutdown timeout (allow MCP server to finish)
# Emergency override: sudo systemctl kill --signal=SIGKILL plant-care-mcp.service
TimeoutStopSec=3700

# Security hardening
NoNewPrivileges=true
PrivateTmp=true


[Install]
WantedBy=multi-user.target
