# app/deploy/plant-care-mcp.service
[Unit]
Description=Plant Care MCP Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mcpserver
Group=mcpserver
WorkingDirectory=/home/mcpserver/plant-care-app
Environment="PATH=/home/mcpserver/.local/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/mcpserver/.env
ExecStart=/home/mcpserver/.local/bin/uv run python run_http.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Graceful shutdown timeout (allow agent to finish)
# Wait indefinitely for agent to complete - never interrupt a running agent
# Emergency override: sudo systemctl kill --signal=SIGKILL plant-care-mcp.service
TimeoutStopSec=3700


[Install]
WantedBy=multi-user.target
